# salt/graphite/carbon_supervisor.conf
[program:carbon]
command=/usr/bin/python /opt/graphite/bin/carbon-cache.py --config /etc/graphite/carbon.conf start
autostart=true
autorestart=true
exitcodes=0
