# salt/graphite/riemann_supervisor.conf
[program:riemann]
command=/usr/bin/riemann /etc/riemann/riemann.config
autostart=true
autorestart=true
exitcodes=0
