# salt/publet/rqdashboard.conf
[program:rqdashboard]
directory=/opt/publet/apps/publet
user=deploy
command=/opt/publet/venvs/publet/bin/rq-dashboard --port=5678
autostart=true
autorestart=true
exitcodes=0
killasgroup=true